# systemd/train-display-reboot.target.service
[Unit]
Description=Trigger weekly system reboot for Train Departure Display
# REBOOT-04: reboot performed via systemctl reboot (clean shutdown sequence)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot
