# systemd/train-display-reboot.timer
[Unit]
Description=Weekly reboot for Train Departure Display
Requires=train-display-reboot.target.service

[Timer]
# Default: every Sunday at 03:00 local time
# Change via install.sh or edit OnCalendar directly
OnCalendar=Sun *-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
